# shoper/config/templates/systemd/frappe-shoper-redis.target
[Unit]
After=network.target
Wants={{ shoper_name }}-redis-cache.service {{ shoper_name }}-redis-queue.service

[Install]
WantedBy=multi-user.target
