# shoper/config/templates/systemd/frappe-shoper-web.target
[Unit]
After=network.target
Wants={{ shoper_name }}-shprho-web.service {{ shoper_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
